# backend/deploy/azure/voiceai-frontend.service
[Unit]
Description=Voice AI Assistant frontend (Next.js, native — no Docker)
After=network-online.target voiceai-backend.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__APP_DIR__/frontend
ExecStart=/usr/bin/npm run start -- -p 3000
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
